# demohouse/deep_search_mcp/backend/supervisord.d/link-search-mcp.conf
[program:link-search-mcp]
command=mcp-server-ark -t sse
directory=.                ; directory to cwd to before exec (def no cwd)
; startsecs=1                   ; # of secs prog must stay up to be running (def. 1)
autorestart=true        ; when to restart if exited after running (def: unexpected)
startsecs=10
stdout_logfile=/tmp/link-search-mcp.std.out
stdout_logfile_maxbytes=32MB
stderr_logfile=/tmp/link-search-mcp.err.out
stderr_logfile_maxbytes=32MB
environment=PORT="7001",ARK_BOT_NAME="web_search",ARK_BOT_DESCRIPTION="这是联网搜索工具，如果需要搜索互联网上的内容，请使用此工具。输入为关键词，每次最多一个关键词",ARK_TOOL_LINK_READER="True"