# demohouse/deep_search_mcp/backend/supervisord.d/sandbox-mcp.conf
[program:sandbox-mcp]
command=mcp-server-vefaas-sandbox -t sse
directory=.                ; directory to cwd to before exec (def no cwd)
; startsecs=1                   ; # of secs prog must stay up to be running (def. 1)
autorestart=true        ; when to restart if exited after running (def: unexpected)
startsecs=180
stdout_logfile=/tmp/sandbox.std.out
stdout_logfile_maxbytes=32MB
stderr_logfile=/tmp/sandbox.err.out
stderr_logfile_maxbytes=32MB
environment=PORT="7002"
