# demohouse/deep_search_mcp/backend/supervisord.d/vdb-mcp.conf
[program:knowledgebase-mcp]
command=mcp-server-knowledgebase -t sse
directory=.                ; directory to cwd to before exec (def no cwd)
; startsecs=1                   ; # of secs prog must stay up to be running (def. 1)
autorestart=true        ; when to restart if exited after running (def: unexpected)
startsecs=10
stdout_logfile=/tmp/vdb.std.out
stdout_logfile_maxbytes=32MB
stderr_logfile=/tmp/vdb.err.out
stderr_logfile_maxbytes=32MB
environment=PORT="7004"